# design_week/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:design_week]
directory=/usr/design_week
command=gunicorn -b 0.0.0.0:7000 -t 180 --access-logfile "-" app:app
